# sync-debian-package/cron.weekly
#!/bin/sh
set -e
cd /var/archive-kit/synchrotron/sync-debian-package/

# fetch stuff from Debian and Tanglu archives
./prepare-sync.sh

# sync packages
sync-debian-package -a testing aequorea main
sync-debian-package -a testing aequorea contrib
sync-debian-package -a testing aequorea non-free

# update archive
dak generate-packages-sources2
dak generate-releases
